# configs/etc/s6/insync/run
#!/bin/bash

exec /usr/bin/insync-headless start --no-daemon
