# configs/etc/s6/insync/run
#!/bin/bash

/usr/bin/insync-headless start --no-daemon
